# daemon/com.modelrealignment.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.modelrealignment.daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>__PROJECT_PATH_PLACEHOLDER__/main_loop.py</string>
        <string>--daemon</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>__PROJECT_PATH_PLACEHOLDER__</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>__PROJECT_PATH_PLACEHOLDER__/logs/daemon_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>__PROJECT_PATH_PLACEHOLDER__/logs/daemon_stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>__PROJECT_PATH_PLACEHOLDER__</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>Nice</key>
    <integer>5</integer>
    
    <key>StartInterval</key>
    <integer>3600</integer>
    
    <key>ThrottleInterval</key>
    <integer>300</integer>
    
    <key>ExitTimeOut</key>
    <integer>30</integer>
</dict>
</plist>