# config/com.ryan.modelrealignment.rewards.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ryan.modelrealignment.rewards</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_PATH_PLACEHOLDER__/scripts/schedule_rewards.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>__PROJECT_PATH_PLACEHOLDER__</string>
    
    <key>StartInterval</key>
    <integer>43200</integer> <!-- 12 hours in seconds -->
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>__PROJECT_PATH_PLACEHOLDER__/logs/rewards_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>__PROJECT_PATH_PLACEHOLDER__/logs/rewards_stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>__PROJECT_PATH_PLACEHOLDER__</string>
    </dict>
</dict>
</plist>